# gnome-extra/iio-sensor-proxy/files/iio-sensor-proxy.openrc
#!/sbin/openrc-run

depend() {
	need dbus
}

start () {
	ebegin "Starting iio-sensor-proxy"
	start-stop-daemon --start --exec /usr/libexec/iio-sensor-proxy \
	    --background --pidfile /var/run/iio-sensor-proxy.pid --make-pidfile
	eend $?
}

stop () {
	ebegin "Stopping iio-sensor-proxy"
	start-stop-daemon --stop --exec /usr/libexec/iio-sensor-proxy \
	    --pidfile /var/run/iio-sensor-proxy.pid
	eend $?
}
